# packages/too-nuts/too_nuts-0.3.tar.gz/too_nuts-0.3/scripts/listeners/tns_source_parser.service
[Unit]
Description=ToO source catalogue listener

[Service]
Type=simple
User=tobi
ExecStart=/abs/path/to/python /abs/path/to/TNS_listener.py /abs/path/to/config.yml
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
